# data/services/systemd/kiran-authentication-daemon.service.in
[Unit]
Description = Kiran Authentication Daemon Service

[Service]
Type=dbus
BusName=com.kylinsec.Kiran.Authentication
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kiran-authentication-daemon
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
