# data/kiran-authentication-devices.service.in
[Unit]
Description=Kiran Authentication Devices Service

[Service]
Type=dbus
BusName=com.kylinsec.Kiran.AuthDevice
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kiran-authentication-devices

[Install]
WantedBy=multi-user.target

